# waveshare-epaper-display.timer.example
[Unit]
Description=Waveshare ePaper display update timer

[Timer]
OnActiveSec=0
OnUnitActiveSec=1m
Unit=waveshare-epaper-display.service

[Install]
WantedBy=timers.target
